# deploy/oraclecloud/forex-dashboard-watchdog.timer
[Unit]
Description=Run Forex Dashboard HTTP Watchdog

[Timer]
OnBootSec=45
OnUnitActiveSec=30
AccuracySec=5
Unit=forex-dashboard-watchdog.service

[Install]
WantedBy=timers.target
